# examples/rendered_configs/systemd/mycelium-demo-postgres.service
# PostgreSQL systemd service for mycelium-demo
# Generated: 2025-10-14 17:39:41
# Installation: sudo cp mycelium-demo-postgres.service /etc/systemd/system/
#              sudo systemctl daemon-reload
#              sudo systemctl enable mycelium-demo-postgres.service
#              sudo systemctl start mycelium-demo-postgres.service

[Unit]
Description=PostgreSQL Database Server for mycelium-demo
After=network.target network-online.target
Wants=network-online.target
Documentation=https://www.postgresql.org/docs/

[Service]
Type=notify
User=postgres
Group=postgres

# Environment
Environment=PGDATA=/var/lib/mycelium-demo/postgres
Environment=POSTGRES_DB=mycelium
Environment=POSTGRES_PORT=5432

# Runtime directory
RuntimeDirectory=mycelium-demo-postgres
RuntimeDirectoryMode=0755

# PostgreSQL commands
ExecStartPre=/usr/bin/pg_ctl -D /var/lib/mycelium-demo/postgres initdb || true
ExecStart=/usr/bin/postgres \
    -D /var/lib/mycelium-demo/postgres \
    -p 5432 \
    -c max_connections=100 \
    -c unix_socket_directories=/var/run/postgresql

ExecReload=/bin/kill -HUP $MAINPID
ExecStop=/usr/bin/pg_ctl -D /var/lib/mycelium-demo/postgres stop -m fast

# Security
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/lib/mycelium-demo/postgres
ReadWritePaths=/var/run/postgresql

# Process management
KillMode=mixed
KillSignal=SIGINT
TimeoutSec=0
Restart=always
RestartSec=5s

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mycelium-demo-postgres

# Resource limits
LimitNOFILE=10000
LimitNPROC=10000

[Install]
WantedBy=multi-user.target
